# alice_trading/trading-server.service
[Unit]
Description=Anti-Gravity Institutional Trading Server
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/alice_trading

# Load environment variables from .env file
EnvironmentFile=/home/ubuntu/alice_trading/.env

# Use system Python 3.7 (or venv if you create one)
ExecStart=/usr/bin/python3 /home/ubuntu/alice_trading/server.py

# Auto-restart on any crash, with a 10s delay
Restart=always
RestartSec=10
StartLimitIntervalSec=60
StartLimitBurst=5

# Logging — tail with: journalctl -u trading-server -f
StandardOutput=journal
StandardError=journal
SyslogIdentifier=trading-server

# Resource limits
TimeoutStartSec=30
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
